# docker/zabbix/etc/supervisord.conf
[supervisord]
nodaemon=true
 
[program:httpd]
command=/usr/sbin/httpd -k start

[program:zabbix-server]
command=/usr/sbin/zabbix_server -c /etc/zabbix/zabbix_server.conf
